# rootfs/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash
echo "STARTING SERVICE: Radarr"
UMASK_SET=${UMASK_SET:-000}
umask ${UMASK_SET}

/bin/s6-setuidgid abc /app/Radarr --no-browser -data=/config
